# detectar_intrusos/ids_daemonServer/dpkg-ids-server/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable --now ids-server.service
exit 0
